# rootfs/etc/s6/openssh/run
#!/bin/bash

if [ -f ./setup ]
then
  source ./setup
fi

STARTCMD="/usr/sbin/sshd -E /var/log/sshd.log -D"
exec ${STARTCMD}
