# rootfs/etc/s6/gitea/run
#!/bin/bash

if [ -f ./setup ]
then
  source ./setup
fi

STARTCMD="/usr/bin/gosu git /usr/bin/gitea web -c /etc/gitea/gitea.ini"
exec ${STARTCMD}
